# deploy/ai-zw.service
[Unit]
Description=AI-ZW 电商图像视频生成工具
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-peacture-video/backend
ExecStart=/root/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 6068
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
